# scripts/com.smoltorrent.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.smoltorrent.api</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Replace with the output of: which python3 (inside your venv) -->
        <string>/Users/yuvrajsingh1/Desktop/smoltorrent/.venv/bin/python</string>
        <string>/Users/yuvrajsingh1/Desktop/smoltorrent/backend/api.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/yuvrajsingh1/Desktop/smoltorrent</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/smoltorrent-api.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/smoltorrent-api.log</string>
</dict>
</plist>
